# wwd.timer
[Unit]
Description=Download a Wallhaven wallpaper every 2 hours

[Timer]
OnBootSec=0s
OnActiveSec=2h

[Install]
WantedBy=timers.target
