# wwd.service
[Unit]
Description=Wallhaven wallpaper downloader
Requisite=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/wwd

[Install]
